# mps/home/etc/trillian.cron
# $HopeName: HOMEmm!etc:trillian.cron(trunk.13) $
# minutes hours day-of-month month day-of-week command
10 21  *  *  *  sh -c ". etc/shenv && overnight checkout"
20 21  *  *  *  sh -c "./bin/tidy"
30 21  *  *  *  sh -c ". etc/shenv && overnight ./bin/build/MM_src_trunk sos8gc"
 0  1  *  *  *  sh -c ". etc/shenv && overnight ./bin/build/MM_dylan_honeybee sos8gc"
30  4  *  *  *  sh -c ". etc/shenv && overnight ./bin/build/MM_dylan_cottonwood sos8gc"
